# WinSxS/amd64_userexperience-desktop_31bf3856ad364e35_10.0.22621.1_none_8debb80c5e085dd8/CBS/DesktopStickerEditor/Task.xml
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>Microsoft Corporation</Author>
    <Source>Microsoft Corporation</Source>
    <Description>Auto-composite stickers onto new background image that has been recently changed.</Description>
  </RegistrationInfo>
  <Triggers>
    <!--WNF_SHEL_BACKGROUND_IMAGE_CHANGED = {0xa3bfb875, 0x0d83063e}; -->
    <WnfStateChangeTrigger>
      <StateName>75b8bfa33e06830d</StateName>
    </WnfStateChangeTrigger>
  </Triggers>
  <Principals>
    <Principal>
      <GroupId>S-1-5-4</GroupId>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
  </Settings>
  <Actions>
    <Exec>
      <Command>%LOCALAPPDATA%\Microsoft\WindowsApps\MicrosoftWindows.DesktopStickerEditorCentennial.exe</Command>
      <Arguments>-Silent</Arguments>
    </Exec>
  </Actions>
</Task>